# com.pitchlayer.homes-tracker.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.pitchlayer.homes-tracker</string>

  <key>ProgramArguments</key>
  <array>
    <!-- Update this path to wherever Python 3 lives: `which python3` -->
    <string>/usr/bin/python3</string>
    <string>-m</string>
    <string>homes_tracker</string>
    <string>run</string>
  </array>

  <!-- Full path to your cloned PitchLayer repo -->
  <key>WorkingDirectory</key>
  <string>/Users/YOURUSERNAME/PitchLayer</string>

  <key>EnvironmentVariables</key>
  <dict>
    <key>CONTACT_ID</key>
    <string>REPLACE_WITH_DELISA_PHONE</string>

    <key>GOOGLE_SHEET_ID</key>
    <string>REPLACE_WITH_SHEET_ID</string>

    <key>ANTHROPIC_API_KEY</key>
    <string>REPLACE_WITH_API_KEY</string>

    <key>GOOGLE_CREDENTIALS_FILE</key>
    <string>/Users/YOURUSERNAME/.config/homes_tracker/google_credentials.json</string>
  </dict>

  <!-- Run every 5 minutes -->
  <key>StartInterval</key>
  <integer>300</integer>

  <!-- Also run immediately when loaded -->
  <key>RunAtLoad</key>
  <true/>

  <key>StandardOutPath</key>
  <string>/tmp/homes-tracker.log</string>

  <key>StandardErrorPath</key>
  <string>/tmp/homes-tracker-error.log</string>
</dict>
</plist>
